# deploy/crontab
# 金融新闻报告器定时任务 - Agent系统
# 每天北京时间早上8点执行美股财经新闻查询
0 8 * * * root cd /app && .venv/bin/python scripts/run_agents.py --config config/tasks.yaml >> /var/log/reporter/agents.log 2>&1

# 每天北京时间早上8点10分执行加密货币新闻查询  
10 8 * * * root cd /app && .venv/bin/python scripts/run_agents.py --config config/crypto_tasks.yaml >> /var/log/reporter/crypto_agents.log 2>&1

# 空行（cron 要求） 