# deploy/crontab
# 金融新闻报告器定时任务 - Agent系统
# 每天北京时间下午4点执行所有查询（UTC 时间 08:00）
0 8 * * * root cd /app && .venv/bin/python scripts/run_agents.py >> /var/log/financial-reporter/agents.log 2>&1

# 空行（cron 要求） 